# dockerized-cron/crontab
* * * * * root bash date-script.sh
# Single container cron example
